# edge-chat-service-supervisord.ini
[program:edge-chat-service]
command=/usr/bin/python3 /exdata/edge-chat-service/edge-chat-svr.py --config /exdata/edge-chat-service/config.json --cookies /exdata/edge-chat-service/cookies.json
directory=/exdata/edge-chat-service/
autorestart=true
stderr_logfile=/exdata/edge-chat-service/log/edge-chat-service.err.log
stdout_logfile=/exdata/edge-chat-service/log/edge-chat-service.out.log
environment=
user=root
stopsignal=INT
startsecs=10
